# app.service
[Unit]
Description=My great Go application

[Service]
ExecStart=/home/vagrant/main
Type=simple

[Install]
WantedBy=multi-user.target
